# cookbooks/nova/templates/default/init-nova-novncproxy.conf.erb
description "Nova NoVNC proxy"
author "Vishvananda Ishaya <vishvananda@gmail.com>"

start on runlevel [2345]
stop on runlevel [!2345]

chdir /var/run

pre-start script
   mkdir -p /var/run/nova
   chown nova:root /var/run/nova/

   mkdir -p /var/lock/nova
   chown nova:root /var/lock/nova/
end script

exec start-stop-daemon --start --chuid nova --exec <%= node[:openstack][:install][:source][:path] %>/bin/nova-novncproxy -- --config-file=/etc/nova/nova.conf --web <%= "#{node[:openstack][:install][:source][:path]}/src/novnc" %>