# cookbooks/trove/templates/default/init-trove-conductor.conf.erb
description "Trove conductor server"
author "ccc <ccc@kakaocorp.com>"

start on runlevel [2345]
stop on runlevel [!2345]
#stop on runlevel [016]

respawn

exec start-stop-daemon --start --chuid trove:trove \
   --chdir /var/lib/trove  \
   --exec <%= node[:openstack][:install][:source][:path] %>/bin/trove-conductor \
   -- --config-file=/etc/trove/trove-conductor.conf --logfile=/var/log/trove/trove-conductor.log
