# cookbooks/trove/templates/default/init-trove-guestagent.conf.erb
description "Trove Guest Agent"
author "ccc <ccc@kakaocorp.com>"

start on started mysql
stop on runlevel [!2345]

respawn

exec start-stop-daemon --start \
 --chdir /var/lib/trove  \
 --exec /usr/local/bin/trove-guestagent -- --config-file=/etc/guest_info --config-file=/etc/trove/trove-guestagent.conf --log-dir=/var/log/trove --logfile=guestagent.log --debug
