# cookbooks/swift/templates/default/init-swift-container-auditor.conf.erb
# swift-container-auditor - SWIFT Container Auditor
#
# The swift container auditor.

description     "SWIFT Container Auditor"
author          "Marc Cluet <marc.cluet@ubuntu.com>"

start on runlevel [2345]
stop on runlevel [016]

pre-start script
  if [ -f "/etc/swift/container-server.conf" ]; then
    exec <%= node[:openstack][:install][:source][:path] %>/bin/swift-init container-auditor start
  else
    exit 1
  fi
end script

post-stop exec <%= node[:openstack][:install][:source][:path] %>/bin/swift-init container-auditor stop

