# cookbooks/swift/templates/default/init-swift-proxy.conf.erb
# swift-proxy - SWIFT Proxy Server
#
# The swift proxy server.

description     "SWIFT Proxy Server"
author          "Marc Cluet <marc.cluet@ubuntu.com>"

start on runlevel [2345]
stop on runlevel [016]

pre-start script
  if [ -f "/etc/swift/proxy-server.conf" ]; then
    exec <%= node[:openstack][:install][:source][:path] %>/bin/swift-init proxy-server start
  else
    exit 1
  fi
end script

post-stop exec <%= node[:openstack][:install][:source][:path] %>/bin/swift-init proxy-server stop

