# cookbooks/glance/templates/default/init-glance-registry.conf.erb
description "Glance registry server"
author "Soren Hansen <soren@linux2go.dk>"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec start-stop-daemon --start --chuid glance \
            --chdir /var/lib/glance --name glance-registry \
            --exec <%= node[:openstack][:install][:source][:path] %>/bin/glance-registry
